# http-server.service
[Unit]
Description=Riftech Security System HTTP Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/riftech-cam-security
Environment="PATH=/root/riftech-cam-security/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/riftech-cam-security/venv/bin/python3 /root/riftech-cam-security/http_server.py
Restart=always
RestartSec=10
StandardOutput=append:/root/riftech-cam-security/logs/http.log
StandardError=append:/root/riftech-cam-security/logs/http_error.log

[Install]
WantedBy=multi-user.target
